# net/open5gs/files/upf.init
#!/bin/sh /etc/rc.common

START=80
PROG=/usr/bin/open5gs-upfd

CONFIG=/etc/open5gs/upf.yaml

start() {
	[ -f $CONFIG ] || return 1
	[ -d /var/log/open5gs ] || mkdir -p /var/log/open5gs
	
	(service_start $PROG --config=$CONFIG &> /dev/null )&
}

stop() {
	service_stop $PROG
}

reload() {
	$PROG reload
}


